# etc/cron.d/cleanup-ec2
* * * * * root /home/builder/neutron-thirdparty-ci/bin/cleanup-ec2.sh
