# etc/cron.d/copy-logs
* * * * * root /home/builder/neutron-thirdparty-ci/master/ci-master-copy-logs > /tmp/copy-logs.out 2> /tmp/copy-logs.err
