# etc/cron.d/zuul-watch
*/2 * * * * root /home/builder/neutron-thirdparty-ci/bin/auto_redo.rb
